# etc/rc.d/rc.local
#!/bin/bash

/home/dicom/startup.bash start ALL
